# supervisor/supervisord-user.conf
[inet_http_server]
port=127.0.0.1:29011

[supervisorctl]
serverurl=http://127.0.0.1:29011

[supervisord]
logfile=/dev/shm/supervisord_user.log
nodaemon=false
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:zo-dashboard]
command=bun /home/workspace/Skills/zo-dashboard/scripts/server.ts
directory=/home/workspace
environment=
autostart=true
autorestart=true
stopsignal=TERM
stopasgroup=true
startretries=20
startsecs=3
stdout_logfile=/dev/shm/zo-dashboard.log
stderr_logfile=/dev/shm/zo-dashboard_err.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=4
killasgroup=true

; Sync worker for 24/7 cloud dashboard.
; Requires ZO_DASHBOARD_CLOUD_URL + ZO_SYNC_SECRET env vars in the user's shell profile.
[program:zo-sync]
command=python3 /home/workspace/Skills/zo-dashboard/scripts/sync_to_cloud.py
directory=/home/workspace
autostart=true
autorestart=true
stopsignal=TERM
stopasgroup=true
startretries=20
startsecs=3
stdout_logfile=/dev/shm/zo-sync.log
stderr_logfile=/dev/shm/zo-sync_err.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=4
killasgroup=true

